# wallpaper-engine-daemon.service
[Unit]
Description=Service to change wallpaper every x time.

[Service]
Type=simple
ExecStart=/home/atreus/wallpaper_engine/daemon

[Install]
WantedBy=default.target
